# debian/postinst
ln -sf /usr/lib/azure-proxy-agent/Package_pkgversion /usr/lib/azure-proxy-agent/package
ln -sf /usr/lib/azure-proxy-agent/package/GuestProxyAgent /usr/sbin/azure-proxy-agent

systemctl daemon-reload
systemctl start GuestProxyAgent.service
systemctl enable GuestProxyAgent.service

#DEBHELPER#
