# pkg_debian/postinst

systemctl unmark azure-proxy-agent.service
systemctl daemon-reload
systemctl start azure-proxy-agent.service
systemctl enable azure-proxy-agent.service

#DEBHELPER#
